# .config/systemd/user/mailsync.service
[Unit]
Description=Mailsync
OnFailure=systemd-notify@%n

[Service]
Type=oneshot
ExecStart=mbsync -a
